# config/systemd/user/emacs.service.d/override.conf
[Unit]
Wants=app-org.fcitx.Fcitx5-autostart.service
After=app-org.fcitx.Fcitx5-autostart.service

[Service]
Environment=LC_CTYPE=zh_CN.UTF-8 XMODIFIERS="@im=fcitx5"
