# config/systemd/user/aria2.service
[Unit]
Description=aria2c Service

[Service]
ExecStart=/usr/bin/aria2c --conf-path=/home/erystasius/.config/aria2/daemon.conf
