# config/systemd/user/cursor-fix.service
[Unit]
Description=Fix cursor shape.

[Service]
Type=oneshot
ExecStart=/usr/bin/xsetroot -cursor_name left_ptr

