# config/systemd/user/jgmenu.service
[Unit]
Description=Start jgmenu daemon.
Before=sxhkd.service

[Service]
ExecStart=/usr/bin/jgmenu_run
