# config/systemd/user/print-env.service
[Unit]
Description=Print environment variables.

[Service]
Type=oneshot
Environment="HHAHAHAH=asfdasfda"
ExecStart=/usr/bin/env

