# config/systemd/user/sxhkd.service
[Unit]
Description=Start picom the compositor.

[Service]
ExecStart=/usr/bin/sxhkd

