# config/systemd/machine.service
[Unit]
Description=Configuring for different machine.

[Service]
Type=oneshot
User=root
Group=root
ExecStart=/usr/bin/whoami

