# config/systemd/polybar@.service
[Unit]
Description=Start polybar.

[Service]
ExecStart=/usr/bin/polybar -r %i

