# config/systemd/services.target
[Unit]
Description=Start my favorite services.

Wants=syncthing.service
Wants=clash.service

[Install]
WantedBy=default.target
