# config/systemd/wallpaper.service
[Unit]
Description=Reset wallpaper with nitrogen.

[Service]
Type=oneshot
ExecStart=/usr/bin/nitrogen --restore

